# hypr-power.service
[Unit]
Description=hypr-power: Lid / Monitor / Lock / Suspend state machine
After=graphical-session.target
PartOf=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hypr-power daemon
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
